# evcc-raedian/rootfs/etc/services.d/evcc/run
#!/usr/bin/with-contenv bashio
exec /usr/local/bin/evcc --config /config/evcc.yaml
